# target-image/ctf/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
chown=root:maintops

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:web]
command=python3 /app/app.py
directory=/app
user=webuser
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/var/log/supervisor/web.log
stderr_logfile=/var/log/supervisor/web_err.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB

[program:sshd]
command=/usr/sbin/sshd -D -e
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd_err.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
